# systemd/notes.service
[Unit]
Description=Auto push notes

[Service]
Type=oneshot
ExecStart=/home/george/.config/scripts/pushnotes.sh
